# install/etc/profile.d/pathenv.sh
export PATH=/var/www/html/bin:/var/www/html/vendor/bin:/var/www/html/vendor/drush/drush:$PATH
